# templates/dropbox.service
[Unit]
Description=Dropbox

[Service]
ExecStart=/usr/bin/dropbox start -i
Type=forking
PIDFile=%h/.dropbox/dropbox.pid

[Install]
WantedBy=default.target
